# debian/usensord.conf
description "usensord service for haptics support"

start on started dbus
stop on stopped dbus

exec usensord
