# debian/usensord.conf
description		"Starts the usensord service which provides haptics support"

start on dbus
stop on runlevel [06]

exec usensord
